# init/voidrift.service
[Unit]
Description=Void Drift IRC idle RPG bot
Documentation=https://github.com/cstroie/voidrift
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=voidrift
Group=voidrift

# Environment file — copy voidrift.env.example to /etc/voidrift/voidrift.env
# and set your values there. The file must be root:root 0600.
# Command-line flags in ExecStart override anything set here.
EnvironmentFile=-/etc/voidrift/voidrift.env

WorkingDirectory=/var/lib/voidrift

ExecStart=/usr/local/bin/voidrift -data voidrift.json -guilds guilds.json -log /var/log/voidrift/voidrift.log

Restart=on-failure
RestartSec=30s

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/voidrift /var/log/voidrift
LogsDirectory=voidrift
LogsDirectoryMode=0750
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
